# run/systemd/generator/sendsigs.service.d/50-insserv.conf-$remote_fs.conf
# Automatically generated by systemd-insserv-generator

[Unit]
Wants=remote-fs-pre.target
Before=remote-fs-pre.target
