# run/systemd/generator/umountfs.service.d/50-insserv.conf-$local_fs.conf
# Automatically generated by systemd-insserv-generator

[Unit]
Wants=local-fs.target
Before=local-fs.target
