# var/lib/dpkg/info/asterisk.conffiles
/etc/default/asterisk
/etc/init.d/asterisk
/etc/logrotate.d/asterisk
