# var/lib/dpkg/info/asterisk.list
/.
/etc
/etc/default
/etc/default/asterisk
/etc/asterisk
/etc/init.d
/etc/init.d/asterisk
/etc/logrotate.d
/etc/logrotate.d/asterisk
/usr
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/asterisk-config-custom.1.gz
/usr/share/man/man1/astman.1.gz
/usr/share/man/man8
/usr/share/man/man8/astdb2bdb.8.gz
/usr/share/man/man8/safe_asterisk.8.gz
/usr/share/man/man8/smsq.8.gz
/usr/share/man/man8/asterisk.8.gz
/usr/share/man/man8/autosupport.8.gz
/usr/share/man/man8/astgenkey.8.gz
/usr/share/man/man8/astdb2sqlite3.8.gz
/usr/share/man/man8/aelparse.8.gz
/usr/share/asterisk
/usr/share/asterisk/agi-bin
/usr/share/asterisk/firmware
/usr/share/asterisk/firmware/iax
/usr/share/asterisk/conf
/usr/share/asterisk/conf/samples
/usr/share/asterisk/conf/samples/res_fax.conf.sample
/usr/share/asterisk/conf/samples/cdr_pgsql.conf.sample
/usr/share/asterisk/conf/samples/func_odbc.conf.sample
/usr/share/asterisk/conf/samples/cel_odbc.conf.sample
/usr/share/asterisk/conf/samples/chan_dahdi.conf.sample
/usr/share/asterisk/conf/samples/pjsip_notify.conf.sample
/usr/share/asterisk/conf/samples/festival.conf.sample
/usr/share/asterisk/conf/samples/console.conf.sample
/usr/share/asterisk/conf/samples/calendar.conf.sample
/usr/share/asterisk/conf/samples/cdr_adaptive_odbc.conf.sample
/usr/share/asterisk/conf/samples/phoneprov.conf.sample
/usr/share/asterisk/conf/samples/vpb.conf.sample
/usr/share/asterisk/conf/samples/cdr_mysql.conf.sample
/usr/share/asterisk/conf/samples/cdr_sqlite3_custom.conf.sample
/usr/share/asterisk/conf/samples/res_pktccops.conf.sample
/usr/share/asterisk/conf/samples/ari.conf.sample
/usr/share/asterisk/conf/samples/cdr_syslog.conf.sample
/usr/share/asterisk/conf/samples/enum.conf.sample
/usr/share/asterisk/conf/samples/ss7.timers.sample
/usr/share/asterisk/conf/samples/test_sorcery.conf.sample
/usr/share/asterisk/conf/samples/agents.conf.sample
/usr/share/asterisk/conf/samples/app_skel.conf.sample
/usr/share/asterisk/conf/samples/res_ldap.conf.sample
/usr/share/asterisk/conf/samples/say.conf.sample
/usr/share/asterisk/conf/samples/ccss.conf.sample
/usr/share/asterisk/conf/samples/asterisk.conf.sample
/usr/share/asterisk/conf/samples/logger.conf.sample
/usr/share/asterisk/conf/samples/cel_pgsql.conf.sample
/usr/share/asterisk/conf/samples/unistim.conf.sample
/usr/share/asterisk/conf/samples/res_stun_monitor.conf.sample
/usr/share/asterisk/conf/samples/statsd.conf.sample
/usr/share/asterisk/conf/samples/adsi.conf.sample
/usr/share/asterisk/conf/samples/res_parking.conf.sample
/usr/share/asterisk/conf/samples/cdr_tds.conf.sample
/usr/share/asterisk/conf/samples/queuerules.conf.sample
/usr/share/asterisk/conf/samples/cdr.conf.sample
/usr/share/asterisk/conf/samples/hep.conf.sample
/usr/share/asterisk/conf/samples/sorcery.conf.sample
/usr/share/asterisk/conf/samples/res_odbc.conf.sample
/usr/share/asterisk/conf/samples/udptl.conf.sample
/usr/share/asterisk/conf/samples/misdn.conf.sample
/usr/share/asterisk/conf/samples/sla.conf.sample
/usr/share/asterisk/conf/samples/http.conf.sample
/usr/share/asterisk/conf/samples/res_config_sqlite3.conf.sample
/usr/share/asterisk/conf/samples/dnsmgr.conf.sample
/usr/share/asterisk/conf/samples/cdr_manager.conf.sample
/usr/share/asterisk/conf/samples/extensions.ael.sample
/usr/share/asterisk/conf/samples/amd.conf.sample
/usr/share/asterisk/conf/samples/oss.conf.sample
/usr/share/asterisk/conf/samples/stasis.conf.sample
/usr/share/asterisk/conf/samples/dundi.conf.sample
/usr/share/asterisk/conf/samples/chan_mobile.conf.sample
/usr/share/asterisk/conf/samples/ooh323.conf.sample
/usr/share/asterisk/conf/samples/extensions.lua.sample
/usr/share/asterisk/conf/samples/smdi.conf.sample
/usr/share/asterisk/conf/samples/motif.conf.sample
/usr/share/asterisk/conf/samples/indications.conf.sample
/usr/share/asterisk/conf/samples/cli_permissions.conf.sample
/usr/share/asterisk/conf/samples/telcordia-1.adsi
/usr/share/asterisk/conf/samples/cdr_custom.conf.sample
/usr/share/asterisk/conf/samples/iaxprov.conf.sample
/usr/share/asterisk/conf/samples/rtp.conf.sample
/usr/share/asterisk/conf/samples/res_snmp.conf.sample
/usr/share/asterisk/conf/samples/mgcp.conf.sample
/usr/share/asterisk/conf/samples/voicemail.conf.sample
/usr/share/asterisk/conf/samples/meetme.conf.sample
/usr/share/asterisk/conf/samples/modules.conf.sample
/usr/share/asterisk/conf/samples/confbridge.conf.sample
/usr/share/asterisk/conf/samples/extensions.conf.sample
/usr/share/asterisk/conf/samples/xmpp.conf.sample
/usr/share/asterisk/conf/samples/cli_aliases.conf.sample
/usr/share/asterisk/conf/samples/manager.conf.sample
/usr/share/asterisk/conf/samples/codecs.conf.sample
/usr/share/asterisk/conf/samples/skinny.conf.sample
/usr/share/asterisk/conf/samples/cli.conf.sample
/usr/share/asterisk/conf/samples/res_curl.conf.sample
/usr/share/asterisk/conf/samples/dbsep.conf.sample
/usr/share/asterisk/conf/samples/res_pgsql.conf.sample
/usr/share/asterisk/conf/samples/res_config_sqlite.conf.sample
/usr/share/asterisk/conf/samples/acl.conf.sample
/usr/share/asterisk/conf/samples/queues.conf.sample
/usr/share/asterisk/conf/samples/cel_tds.conf.sample
/usr/share/asterisk/conf/samples/followme.conf.sample
/usr/share/asterisk/conf/samples/res_config_mysql.conf.sample
/usr/share/asterisk/conf/samples/dsp.conf.sample
/usr/share/asterisk/conf/samples/config_test.conf.sample
/usr/share/asterisk/conf/samples/sip_notify.conf.sample
/usr/share/asterisk/conf/samples/extconfig.conf.sample
/usr/share/asterisk/conf/samples/pjsip.conf.sample
/usr/share/asterisk/conf/samples/cel_custom.conf.sample
/usr/share/asterisk/conf/samples/asterisk.adsi
/usr/share/asterisk/conf/samples/sip.conf.sample
/usr/share/asterisk/conf/samples/alarmreceiver.conf.sample
/usr/share/asterisk/conf/samples/users.conf.sample
/usr/share/asterisk/conf/samples/cel.conf.sample
/usr/share/asterisk/conf/samples/extensions_minivm.conf.sample
/usr/share/asterisk/conf/samples/cel_sqlite3_custom.conf.sample
/usr/share/asterisk/conf/samples/osp.conf.sample
/usr/share/asterisk/conf/samples/res_corosync.conf.sample
/usr/share/asterisk/conf/samples/minivm.conf.sample
/usr/share/asterisk/conf/samples/muted.conf.sample
/usr/share/asterisk/conf/samples/cdr_odbc.conf.sample
/usr/share/asterisk/conf/samples/features.conf.sample
/usr/share/asterisk/conf/samples/alsa.conf.sample
/usr/share/asterisk/conf/samples/iax.conf.sample
/usr/share/asterisk/conf/samples/app_mysql.conf.sample
/usr/share/asterisk/conf/samples/phone.conf.sample
/usr/share/asterisk/conf/samples/musiconhold.conf.sample
/usr/share/asterisk/keys
/usr/share/asterisk/documentation
/usr/share/asterisk/documentation/appdocsxml.dtd
/usr/share/asterisk/documentation/appdocsxml.xslt
/usr/share/asterisk/documentation/core-en_US.xml
/usr/share/asterisk/documentation/thirdparty
/usr/share/asterisk/static-http
/usr/share/asterisk/sounds
/usr/share/doc
/usr/share/doc/asterisk
/usr/share/doc/asterisk/CREDITS.gz
/usr/share/doc/asterisk/README.Debian.gz
/usr/share/doc/asterisk/menuselect.makeopts.gz
/usr/share/doc/asterisk/README.gz
/usr/share/doc/asterisk/examples
/usr/share/doc/asterisk/examples/asterisk_cleanup
/usr/share/doc/asterisk/examples/sample.call
/usr/share/doc/asterisk/examples/vmail.cgi.gz
/usr/share/doc/asterisk/examples/asterisk.startup.conf
/usr/share/doc/asterisk/BUGS
/usr/share/doc/asterisk/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/asterisk
/usr/sbin
/usr/sbin/astdb2bdb
/usr/sbin/safe_asterisk
/usr/sbin/astdb2sqlite3
/usr/sbin/autosupport
/usr/sbin/asterisk
/usr/sbin/aelparse
/usr/sbin/astgenkey
/usr/sbin/smsq
/usr/sbin/astcanary
/usr/bin
/usr/bin/asterisk-config-custom
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/asterisk.service
/var
/var/spool
/var/spool/asterisk
/var/spool/asterisk/system
/var/spool/asterisk/dictate
/var/spool/asterisk/monitor
/var/spool/asterisk/tmp
/var/spool/asterisk/outgoing
/var/spool/asterisk/meetme
/var/spool/asterisk/voicemail
/var/lib
/var/lib/asterisk
/var/lib/asterisk/sounds
/var/lib/asterisk/sounds/custom
/var/lib/asterisk/moh
/var/log
/var/log/asterisk
/var/log/asterisk/cdr-csv
/var/log/asterisk/cdr-custom
/usr/share/man/man8/rasterisk.8.gz
/usr/share/asterisk/sounds/recordings
/usr/share/asterisk/sounds/custom
/usr/share/doc/asterisk/NEWS.Debian.gz
/usr/share/doc/asterisk/changelog.Debian.gz
/usr/sbin/rasterisk
