# var/lib/dpkg/info/landscape-client.conffiles
/etc/dbus-1/system.d/landscape.conf
/etc/init.d/landscape-client
/etc/logrotate.d/landscape-client
