# var/lib/dpkg/info/landscape-client.list
/.
/etc
/etc/logrotate.d
/etc/logrotate.d/landscape-client
/etc/init.d
/etc/init.d/landscape-client
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/landscape.conf
/usr
/usr/bin
/usr/bin/landscape-monitor
/usr/bin/landscape-config
/usr/bin/landscape-message
/usr/bin/landscape-manager
/usr/bin/landscape-broker
/usr/bin/landscape-package-changer
/usr/bin/landscape-dbus-proxy
/usr/bin/landscape-client
/usr/bin/landscape-package-reporter
/usr/bin/landscape-release-upgrader
/usr/lib
/usr/lib/landscape
/usr/lib/landscape/apt-update
/usr/share
/usr/share/doc
/usr/share/doc/landscape-client
/usr/share/doc/landscape-client/copyright
/usr/share/doc/landscape-client/example.conf.gz
/usr/share/doc/landscape-client/README
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/landscape-client.1.gz
/usr/share/man/man1/landscape-message.1.gz
/usr/share/man/man1/landscape-config.1.gz
/usr/share/doc/landscape-client/changelog.Debian.gz
