# var/lib/dpkg/info/lxd.conffiles
/etc/apparmor.d/usr.lib.lxd.lxd-bridge-proxy
/etc/init.d/lxd
/etc/init/lxd.conf
