# var/lib/dpkg/info/php-common.list
/.
/etc
/etc/cron.d
/etc/cron.d/php
/var
/var/lib
/var/lib/php
/var/lib/php/modules
/var/lib/php/sessions
/usr
/usr/sbin
/usr/sbin/phpquery
/usr/sbin/phpenmod
/usr/lib
/usr/lib/php
/usr/lib/php/sessionclean
/usr/lib/php/php-maintscript-helper
/usr/lib/php/php-helper
/usr/share
/usr/share/doc
/usr/share/doc/php-common
/usr/share/doc/php-common/changelog.gz
/usr/share/doc/php-common/copyright
/usr/sbin/phpdismod
