# var/lib/dpkg/info/php7.0-fpm.conffiles
/etc/apache2/conf-available/php7.0-fpm.conf
/etc/init.d/php7.0-fpm
/etc/init/php7.0-fpm.conf
/etc/logrotate.d/php7.0-fpm
/etc/php/7.0/fpm/php-fpm.conf
/etc/php/7.0/fpm/pool.d/www.conf
