# run/systemd/generator/hwclock.service.d/50-insserv.conf-$time.conf
# Automatically generated by systemd-insserv-generator

[Unit]
Wants=time-sync.target
Before=time-sync.target
